# src/systemd/root/zram-init.service
[Unit]
Description=Initialise zram.

[Service]
Type=simple
ExecStart=BIN_PREFIX/zram-init

[Install]
WantedBy=default.target
