# src/systemd/root/zram-init.service
[Unit]
Description=Initialise zram.

[Service]
Type=simple
ExecStart=/usr/local/bin/zram-init

[Install]
WantedBy=default.target
