# src/systemd/user/logboot.service
[Unit]
Description=Log boot analysis.
After=multi-user.target

[Service]
Type=simple
ExecStart=BIN_PREFIX/logboot

[Install]
WantedBy=default.target
